# vex-daemon.service
[Unit]
Description=Vex Daemon — identity continuity bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/aldous/Desktop/vex
Environment=VEX_HOST=0.0.0.0
ExecStart=/home/aldous/Desktop/vex/.venv/bin/python3 -m vex_daemon.daemon
Restart=always
RestartSec=5
User=aldous

[Install]
WantedBy=default.target
